# deploy/fortios-upgrade.service
[Unit]
Description=FortiOS Upgrade Intelligence
After=network.target

[Service]
Type=simple
User=tetrax
Group=tetrax
WorkingDirectory=/home/tetrax/workspace/upgrade_path
ExecStart=/usr/bin/python3 scripts/fortios_server.py --host 127.0.0.1 --port 8000
Restart=on-failure
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/tetrax/workspace/upgrade_path/data /home/tetrax/workspace/upgrade_path/docs

[Install]
WantedBy=multi-user.target
