# deploy/fortios-catalog-refresh.timer
[Unit]
Description=Daily refresh of the FortiOS Upgrade Intelligence catalog

[Timer]
OnCalendar=*-*-* 07:15:00
Persistent=true

[Install]
WantedBy=timers.target
